# project/dist/play.conf
# Updstart configuration file

description "Play Server"

start on runlevel [2345]
stop on runlevel [12345]

setuid playd

console log

env PLAY_DIR=/home/ubuntu/soq/current


exec /usr/bin/java -Dhttp.port=9001 -classpath $PLAY_DIR/lib\* play.core.server.NettyServer $PLAY_DIR